# services/host/ros.service
[Unit]
After=network.target

[Service]
Type=simple
ExecStart=/home/astrib/bees/scripts/start/host/start_ros.bash

[Install]
WantedBy=default.target