# start/host.service
[Unit]
After=network.target

[Service]
Type=simple
ExecStart=/home/astrib/bees/scripts/start/start_host.bash

[Install]
WantedBy=default.target