# start/robot.service
[Unit]
After=network.target

[Service]
Type=simple
ExecStart=/home/hiver/bees/scripts/start/start_robot.bash

[Install]
WantedBy=default.target